# systemd/vm-pcloud-sync.timer
[Unit]
Description=Run vm-pcloud-sync every 15 minutes

[Timer]
OnBootSec=2m
OnUnitActiveSec=15m
Persistent=true

[Install]
WantedBy=timers.target
